# api-node/supervisor/aen-server.conf
[program:aen-server]
command=/etc/init.d/aen-chain-startup
priority=900
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
